# upstart/demo.conf
description "demo"

start on (filesystem)
stop on runlevel [016]

respawn
console log
setuid nobody
setgid nogroup
chdir /var/www/demo

exec env/bin/gunicorn --config gunicorn.conf demo:app
